# root/usr/lib/systemd/system/slapd-create-dhparam.service
[Unit]
Description=Create the dhparam group for slapd

[Service]
Type=simple
ExecStart=/usr/libexec/nethserver/slapd-create-dhparam
