# root/etc/systemd/system/slapd.service.d/nethserver.conf
[Unit]
Wants=slapd-create-dhparam.service
